# var/lib/dpkg/info/cluster-glue.list
/.
/var
/var/lib
/var/lib/heartbeat
/var/lib/heartbeat/lrm
/var/lib/heartbeat/cores
/var/lib/heartbeat/cores/hacluster
/var/lib/heartbeat/cores/root
/var/lib/heartbeat/cores/nobody
/usr
/usr/lib
/usr/lib/stonith
/usr/lib/stonith/plugins
/usr/lib/stonith/plugins/xen0-ha-dom0-stonith-helper
/usr/lib/stonith/plugins/stonith2
/usr/lib/stonith/plugins/stonith2/ssh.so
/usr/lib/stonith/plugins/stonith2/nw_rpc100s.so
/usr/lib/stonith/plugins/stonith2/wti_mpc.so
/usr/lib/stonith/plugins/stonith2/drac3.so
/usr/lib/stonith/plugins/stonith2/baytech.so
/usr/lib/stonith/plugins/stonith2/apcmaster.so
/usr/lib/stonith/plugins/stonith2/suicide.so
/usr/lib/stonith/plugins/stonith2/external.so
/usr/lib/stonith/plugins/stonith2/ribcl.py
/usr/lib/stonith/plugins/stonith2/cyclades.so
/usr/lib/stonith/plugins/stonith2/rcd_serial.so
/usr/lib/stonith/plugins/stonith2/rps10.so
/usr/lib/stonith/plugins/stonith2/bladehpi.so
/usr/lib/stonith/plugins/stonith2/null.so
/usr/lib/stonith/plugins/stonith2/rhcs.so
/usr/lib/stonith/plugins/stonith2/apcmastersnmp.so
/usr/lib/stonith/plugins/stonith2/wti_nps.so
/usr/lib/stonith/plugins/stonith2/apcsmart.so
/usr/lib/stonith/plugins/stonith2/ibmhmc.so
/usr/lib/stonith/plugins/stonith2/meatware.so
/usr/lib/stonith/plugins/stonith2/ipmilan.so
/usr/lib/stonith/plugins/external
/usr/lib/stonith/plugins/external/ippower9258
/usr/lib/stonith/plugins/external/riloe
/usr/lib/stonith/plugins/external/xen0-ha
/usr/lib/stonith/plugins/external/rackpdu
/usr/lib/stonith/plugins/external/sbd
/usr/lib/stonith/plugins/external/dracmc-telnet
/usr/lib/stonith/plugins/external/kdumpcheck
/usr/lib/stonith/plugins/external/vmware
/usr/lib/stonith/plugins/external/hmchttp
/usr/lib/stonith/plugins/external/xen0
/usr/lib/stonith/plugins/external/ssh
/usr/lib/stonith/plugins/external/ipmi
/usr/lib/stonith/plugins/external/drac5
/usr/lib/stonith/plugins/external/ibmrsa-telnet
/usr/lib/stonith/plugins/external/ibmrsa
/usr/lib/heartbeat
/usr/lib/heartbeat/ipctransientserver
/usr/lib/heartbeat/base64_md5_test
/usr/lib/heartbeat/ipctransientclient
/usr/lib/heartbeat/plugins
/usr/lib/heartbeat/plugins/test
/usr/lib/heartbeat/plugins/test/test.so
/usr/lib/heartbeat/plugins/InterfaceMgr
/usr/lib/heartbeat/plugins/InterfaceMgr/generic.so
/usr/lib/heartbeat/plugins/RAExec
/usr/lib/heartbeat/plugins/RAExec/heartbeat.so
/usr/lib/heartbeat/plugins/RAExec/lsb.so
/usr/lib/heartbeat/plugins/RAExec/ocf.so
/usr/lib/heartbeat/ha_logd
/usr/lib/heartbeat/ipctest
/usr/lib/heartbeat/transient-test.sh
/usr/lib/heartbeat/lrmd
/usr/lib/heartbeat/logtest
/usr/share
/usr/share/cluster-glue
/usr/share/cluster-glue/combine-logs.pl
/usr/share/cluster-glue/openais_conf_support.sh
/usr/share/cluster-glue/lrmtest
/usr/share/cluster-glue/lrmtest/lrmregtest-heartbeat
/usr/share/cluster-glue/lrmtest/testcases
/usr/share/cluster-glue/lrmtest/testcases/rscmgmt.exp
/usr/share/cluster-glue/lrmtest/testcases/basicset
/usr/share/cluster-glue/lrmtest/testcases/rscmgmt
/usr/share/cluster-glue/lrmtest/testcases/metadata
/usr/share/cluster-glue/lrmtest/testcases/BSC
/usr/share/cluster-glue/lrmtest/testcases/metadata.exp
/usr/share/cluster-glue/lrmtest/testcases/xmllint.sh
/usr/share/cluster-glue/lrmtest/testcases/common.filter
/usr/share/cluster-glue/lrmtest/testcases/rscexec.exp
/usr/share/cluster-glue/lrmtest/testcases/ra-list.sh
/usr/share/cluster-glue/lrmtest/testcases/rscexec
/usr/share/cluster-glue/lrmtest/testcases/rscmgmt.log_filter
/usr/share/cluster-glue/lrmtest/lrmadmin-interface
/usr/share/cluster-glue/lrmtest/evaltest.sh
/usr/share/cluster-glue/lrmtest/README.regression
/usr/share/cluster-glue/lrmtest/lrmregtest-lsb
/usr/share/cluster-glue/lrmtest/descriptions
/usr/share/cluster-glue/lrmtest/LRMBasicSanityCheck
/usr/share/cluster-glue/lrmtest/defaults
/usr/share/cluster-glue/lrmtest/language
/usr/share/cluster-glue/lrmtest/regression.sh
/usr/share/cluster-glue/lrmtest/lrmregtest
/usr/share/cluster-glue/ha_log.sh
/usr/share/cluster-glue/utillib.sh
/usr/share/cluster-glue/ha_cf_support.sh
/usr/share/doc
/usr/share/doc/cluster-glue
/usr/share/doc/cluster-glue/changelog.Debian.gz
/usr/share/doc/cluster-glue/stonith
/usr/share/doc/cluster-glue/stonith/README.cyclades
/usr/share/doc/cluster-glue/stonith/README.vacm
/usr/share/doc/cluster-glue/stonith/README.external
/usr/share/doc/cluster-glue/stonith/README.drac3
/usr/share/doc/cluster-glue/stonith/README.meatware
/usr/share/doc/cluster-glue/stonith/README.riloe
/usr/share/doc/cluster-glue/stonith/README.bladehpi
/usr/share/doc/cluster-glue/stonith/README.wti_mpc
/usr/share/doc/cluster-glue/stonith/README.dracmc.gz
/usr/share/doc/cluster-glue/stonith/README.ibmrsa-telnet
/usr/share/doc/cluster-glue/stonith/README_kdumpcheck.txt.gz
/usr/share/doc/cluster-glue/stonith/README.rcd_serial.gz
/usr/share/doc/cluster-glue/stonith/README.ibmrsa
/usr/share/doc/cluster-glue/stonith/README.ipmilan.gz
/usr/share/doc/cluster-glue/changelog.gz
/usr/share/doc/cluster-glue/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/stonith.8.gz
/usr/share/man/man8/ha_logd.8.gz
/usr/share/man/man8/hb_report.8.gz
/usr/share/man/man8/meatclient.8.gz
/usr/share/man/man1
/usr/share/man/man1/ha_logger.1.gz
/usr/sbin
/usr/sbin/sbd
/usr/sbin/stonith
/usr/sbin/hb_report
/usr/sbin/lrmadmin
/usr/sbin/ha_logger
/usr/sbin/meatclient
/etc
/etc/init.d
/etc/init.d/logd
