# var/lib/dpkg/info/conntrackd.conffiles
/etc/logrotate.d/conntrackd
/etc/conntrackd/conntrackd.conf
/etc/init.d/conntrackd
/etc/default/conntrackd
