# var/lib/dpkg/info/conntrackd.list
/.
/usr
/usr/sbin
/usr/sbin/conntrackd
/usr/share
/usr/share/doc
/usr/share/doc/conntrackd
/usr/share/doc/conntrackd/copyright
/usr/share/doc/conntrackd/changelog.Debian.gz
/usr/share/doc/conntrackd/AUTHORS
/usr/share/doc/conntrackd/README.Debian
/usr/share/doc/conntrackd/examples
/usr/share/doc/conntrackd/examples/stats
/usr/share/doc/conntrackd/examples/stats/conntrackd.conf
/usr/share/doc/conntrackd/examples/sync
/usr/share/doc/conntrackd/examples/sync/keepalived-multi2.conf
/usr/share/doc/conntrackd/examples/sync/primary-backup.sh
/usr/share/doc/conntrackd/examples/sync/alarm
/usr/share/doc/conntrackd/examples/sync/alarm/conntrackd.conf.gz
/usr/share/doc/conntrackd/examples/sync/alarm/README
/usr/share/doc/conntrackd/examples/sync/notrack
/usr/share/doc/conntrackd/examples/sync/notrack/conntrackd.conf.gz
/usr/share/doc/conntrackd/examples/sync/notrack/README
/usr/share/doc/conntrackd/examples/sync/ftfw
/usr/share/doc/conntrackd/examples/sync/ftfw/conntrackd.conf.gz
/usr/share/doc/conntrackd/examples/sync/ftfw/README
/usr/share/doc/conntrackd/examples/sync/multiprimary.sh.gz
/usr/share/doc/conntrackd/examples/sync/keepalived.conf
/usr/share/doc/conntrackd/examples/sync/keepalived-multi1.conf
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/conntrackd.8.gz
/etc
/etc/logrotate.d
/etc/logrotate.d/conntrackd
/etc/conntrackd
/etc/conntrackd/conntrackd.conf
/etc/init.d
/etc/init.d/conntrackd
/etc/default
/etc/default/conntrackd
