# var/lib/dpkg/info/ddclient.list
/.
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/ddclient.8.gz
/usr/share/doc
/usr/share/doc/ddclient
/usr/share/doc/ddclient/README.cisco
/usr/share/doc/ddclient/changelog.Debian.gz
/usr/share/doc/ddclient/copyright
/usr/share/doc/ddclient/examples
/usr/share/doc/ddclient/examples/sample-etc_rc.d_init.d_ddclient.ubuntu
/usr/share/doc/ddclient/examples/sample-etc_rc.d_init.d_ddclient
/usr/share/doc/ddclient/examples/sample-etc_dhcpc_dhcpcd-eth0.exe
/usr/share/doc/ddclient/examples/sample-etc_rc.d_init.d_ddclient.redhat
/usr/share/doc/ddclient/examples/sample-etc_rc.d_init.d_ddclient.lsb
/usr/share/doc/ddclient/examples/sample-etc_dhclient-exit-hooks
/usr/share/doc/ddclient/examples/sample-etc_ppp_ip-up.local
/usr/share/doc/ddclient/examples/sample-etc_ddclient.conf.gz
/usr/share/doc/ddclient/examples/sample-etc_cron.d_ddclient
/usr/share/doc/ddclient/README.gz
/usr/share/doc/ddclient/README.ssl
/usr/sbin
/usr/sbin/ddclient
/var
/var/cache
/var/cache/ddclient
/etc
/etc/ppp
/etc/ppp/ip-up.d
/etc/ppp/ip-up.d/ddclient
/etc/init.d
/etc/init.d/ddclient
