# var/lib/dpkg/info/ddclient.conffiles
/etc/ppp/ip-up.d/ddclient
/etc/init.d/ddclient
