# var/lib/dpkg/info/eventwatchd.conffiles
/etc/eventwatchd/eventwatchd.conf.sample
/etc/init.d/eventwatchd
