# var/lib/dpkg/info/fuse-utils.list
/.
/usr
/usr/bin
/usr/bin/fusermount
/usr/share
/usr/share/doc
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/fusermount.1.gz
/etc
/etc/init.d
/etc/init.d/fuse
/etc/fuse.conf
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/60-fuse-utils.rules
/sbin
/sbin/mount.fuse
/usr/share/doc/fuse-utils
