# var/lib/dpkg/info/ipsec-tools.list
/.
/etc
/etc/ipsec-tools.conf
/etc/default
/etc/default/setkey
/etc/init.d
/etc/init.d/setkey
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/setkey.8.gz
/usr/share/doc
/usr/share/doc/ipsec-tools
/usr/share/doc/ipsec-tools/copyright
/usr/share/doc/ipsec-tools/changelog.gz
/usr/share/doc/ipsec-tools/README.Debian
/usr/share/doc/ipsec-tools/changelog.Debian.gz
/usr/sbin
/usr/sbin/setkey
/usr/lib
/usr/lib/ipsec-tools
/usr/lib/ipsec-tools/libipsec.so.0.0.1
/usr/lib/ipsec-tools/libipsec.so.0
