# var/lib/dpkg/info/ipsec-tools.conffiles
/etc/ipsec-tools.conf
/etc/default/setkey
/etc/init.d/setkey
