# var/lib/dpkg/info/ipvsadm.list
/.
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/ipvsadm-restore.8.gz
/usr/share/man/man8/ipvsadm-save.8.gz
/usr/share/man/man8/ipvsadm.8.gz
/usr/share/doc
/usr/share/doc/ipvsadm
/usr/share/doc/ipvsadm/SCHEDULERS
/usr/share/doc/ipvsadm/NEWS.Debian.gz
/usr/share/doc/ipvsadm/copyright
/usr/share/doc/ipvsadm/README
/usr/share/doc/ipvsadm/changelog.Debian.gz
/etc
/etc/ipvsadm.rules
/etc/default
/etc/default/ipvsadm
/etc/init.d
/etc/init.d/ipvsadm
/sbin
/sbin/ipvsadm-save
/sbin/ipvsadm
/sbin/ipvsadm-restore
