# var/lib/dpkg/info/ipvsadm.conffiles
/etc/ipvsadm.rules
/etc/default/ipvsadm
/etc/init.d/ipvsadm
