# var/lib/dpkg/info/live-initramfs.list
/.
/usr
/usr/share
/usr/share/live-initramfs
/usr/share/live-initramfs/live-reconfigure
/usr/share/live-initramfs/live-preseed
/usr/share/live-initramfs/languagelist
/usr/share/initramfs-tools
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/live
/usr/share/initramfs-tools/conf.d
/usr/share/initramfs-tools/conf.d/compcache
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/live-premount
/usr/share/initramfs-tools/scripts/live-premount/modules
/usr/share/initramfs-tools/scripts/live-helpers
/usr/share/initramfs-tools/scripts/live
/usr/share/initramfs-tools/scripts/live-functions
/usr/share/initramfs-tools/scripts/live-bottom
/usr/share/initramfs-tools/scripts/live-bottom/34disable_kpersonalizer
/usr/share/initramfs-tools/scripts/live-bottom/15autologin
/usr/share/initramfs-tools/scripts/live-bottom/13swap
/usr/share/initramfs-tools/scripts/live-bottom/10adduser
/usr/share/initramfs-tools/scripts/live-bottom/21xdriver
/usr/share/initramfs-tools/scripts/live-bottom/22screensaver
/usr/share/initramfs-tools/scripts/live-bottom/22gnome_panel_data
/usr/share/initramfs-tools/scripts/live-bottom/02timezone
/usr/share/initramfs-tools/scripts/live-bottom/38jockey_disable_check
/usr/share/initramfs-tools/scripts/live-bottom/08persistence_excludes
/usr/share/initramfs-tools/scripts/live-bottom/23etc_modules
/usr/share/initramfs-tools/scripts/live-bottom/05mountpoints
/usr/share/initramfs-tools/scripts/live-bottom/36disable_trackerd
/usr/share/initramfs-tools/scripts/live-bottom/14locales
/usr/share/initramfs-tools/scripts/live-bottom/35fix_language_selector
/usr/share/initramfs-tools/scripts/live-bottom/21xvidemode
/usr/share/initramfs-tools/scripts/live-bottom/25configure_init
/usr/share/initramfs-tools/scripts/live-bottom/44pk_allow
/usr/share/initramfs-tools/scripts/live-bottom/32disable_hibernation
/usr/share/initramfs-tools/scripts/live-bottom/12fstab
/usr/share/initramfs-tools/scripts/live-bottom/24preseed
/usr/share/initramfs-tools/scripts/live-bottom/19keyboard
/usr/share/initramfs-tools/scripts/live-bottom/16umountfs
/usr/share/initramfs-tools/scripts/live-bottom/33enable_apport_crashes
/usr/share/initramfs-tools/scripts/live-bottom/20xconfig
/usr/share/initramfs-tools/scripts/live-bottom/31disable_update_notifier
/usr/share/initramfs-tools/scripts/live-bottom/42disable_apparmor
/usr/share/initramfs-tools/scripts/live-bottom/99hook
/usr/share/initramfs-tools/scripts/live-bottom/02etc_live_conf
/usr/share/initramfs-tools/scripts/live-bottom/23networking
/usr/share/initramfs-tools/scripts/live-bottom/41apt_cdrom
/usr/share/initramfs-tools/scripts/live-bottom/06hostname
/usr/share/initramfs-tools/scripts/live-bottom/30accessibility
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/live-initramfs
/usr/share/doc
/usr/share/doc/live-initramfs
/usr/share/doc/live-initramfs/copyright
/usr/share/doc/live-initramfs/CREDITS
/usr/share/doc/live-initramfs/parameters.txt
/usr/share/doc/live-initramfs/changelog.Debian.gz
/usr/share/doc/live-initramfs/AUTHORS
/usr/share/doc/live-initramfs/changelog.gz
/usr/share/doc/live-initramfs/changelog.casper.gz
/usr/share/doc/live-initramfs/examples
/usr/share/doc/live-initramfs/examples/live-snapshot.list
/usr/share/doc/live-initramfs/examples/compcache
/usr/share/doc/live-initramfs/examples/live-persistence.binds
/usr/share/doc/live-initramfs/examples/live-snapshot.exclude_list
/usr/share/doc/live-initramfs/examples/live.conf
/usr/share/bug
/usr/share/bug/live-initramfs
/sbin
/sbin/live-snapshot
/sbin/live-login
/sbin/live-getty
/sbin/live-new-uuid
/etc
/etc/init.d
/etc/init.d/live-initramfs
/etc/live.conf
