# var/lib/dpkg/info/live-initramfs.conffiles
/etc/init.d/live-initramfs
/etc/live.conf
