# var/lib/dpkg/info/netplug.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/netplug
/usr/share/doc/netplug/copyright
/usr/share/doc/netplug/NEWS.gz
/usr/share/doc/netplug/changelog.Debian.gz
/usr/share/doc/netplug/README
/usr/share/doc/netplug/TODO
/usr/share/doc/netplug/changelog.gz
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/netplugd.8.gz
/sbin
/sbin/netplugd
/etc
/etc/netplug
/etc/netplug/netplugd.conf
/etc/netplug/netplug
/etc/init.d
/etc/init.d/netplug
