# var/lib/dpkg/info/open-vm-tools.list
/.
/usr
/usr/include
/usr/include/vmGuestLib
/usr/include/vmGuestLib/vm_basic_types.h
/usr/include/vmGuestLib/includeCheck.h
/usr/include/vmGuestLib/vmGuestLib.h
/usr/include/vmGuestLib/vmSessionId.h
/usr/sbin
/usr/share
/usr/share/open-vm-tools
/usr/share/open-vm-tools/messages
/usr/share/open-vm-tools/messages/ja
/usr/share/open-vm-tools/messages/ja/toolboxcmd.vmsg
/usr/share/open-vm-tools/messages/ja/vmtoolsd.vmsg
/usr/share/open-vm-tools/messages/de
/usr/share/open-vm-tools/messages/de/toolboxcmd.vmsg
/usr/share/open-vm-tools/messages/de/vmtoolsd.vmsg
/usr/share/open-vm-tools/messages/zh_CN
/usr/share/open-vm-tools/messages/zh_CN/toolboxcmd.vmsg
/usr/share/open-vm-tools/messages/ko
/usr/share/open-vm-tools/messages/ko/toolboxcmd.vmsg
/usr/share/open-vm-tools/messages/ko/vmtoolsd.vmsg
/usr/share/initramfs-tools
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/vmxnet
/usr/share/doc
/usr/share/doc/open-vm-tools
/usr/share/doc/open-vm-tools/copyright
/usr/share/doc/open-vm-tools/changelog.Debian.gz
/usr/share/doc/open-vm-tools/changelog.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/vmware-xferlogs.1.gz
/usr/share/man/man1/vmware-hgfsclient.1.gz
/usr/share/man/man1/vmware-toolbox-cmd.1.gz
/usr/share/man/man1/vmware-checkvm.1.gz
/usr/share/man/man9
/usr/share/man/man9/vmxnet.9.gz
/usr/share/man/man9/vmci.9.gz
/usr/share/man/man9/vmblock.9.gz
/usr/share/man/man9/vmhgfs.9.gz
/usr/share/man/man9/vmxnet3.9.gz
/usr/share/man/man9/vmmemctl.9.gz
/usr/share/man/man9/vmsync.9.gz
/usr/share/man/man9/vmsock.9.gz
/usr/share/man/man3
/usr/share/man/man3/libguestlib.3.gz
/usr/share/man/man3/libvmtools.3.gz
/usr/share/man/man8
/usr/share/man/man8/vmware-hgfsmounter.8.gz
/usr/lib
/usr/lib/libvmtools.a
/usr/lib/libhgfs.so.0.0.0
/usr/lib/open-vm-tools
/usr/lib/open-vm-tools/plugins
/usr/lib/open-vm-tools/plugins/common
/usr/lib/open-vm-tools/plugins/common/libhgfsServer.so
/usr/lib/open-vm-tools/plugins/common/libvix.so
/usr/lib/open-vm-tools/plugins/vmsvc
/usr/lib/open-vm-tools/plugins/vmsvc/libguestInfo.so
/usr/lib/open-vm-tools/plugins/vmsvc/libpowerOps.so
/usr/lib/open-vm-tools/plugins/vmsvc/libvmbackup.so
/usr/lib/open-vm-tools/plugins/vmsvc/libtimeSync.so
/usr/lib/libguestlib.so.0.0.0
/usr/lib/libguestlib.a
/usr/lib/libvmtools.so.0.0.0
/usr/lib/libhgfs.a
/usr/lib/pkgconfig
/usr/lib/pkgconfig/vmguestlib.pc
/usr/bin
/usr/bin/vmtoolsd
/usr/bin/vmware-rpctool
/usr/bin/vmware-vmblock-fuse
/usr/bin/vmware-hgfsclient
/usr/bin/vmware-xferlogs
/usr/bin/vmware-toolbox-cmd
/usr/bin/vmware-checkvm
/sbin
/sbin/mount.vmhgfs
/etc
/etc/modprobe.d
/etc/modprobe.d/open-vm-tools.conf
/etc/vmware-tools
/etc/vmware-tools/statechange.subr
/etc/vmware-tools/resume-vm-default
/etc/vmware-tools/scripts
/etc/vmware-tools/scripts/vmware
/etc/vmware-tools/scripts/vmware/network
/etc/vmware-tools/vm-support
/etc/vmware-tools/tools.conf
/etc/vmware-tools/suspend-vm-default
/etc/vmware-tools/poweroff-vm-default
/etc/vmware-tools/poweron-vm-default
/etc/pam.d
/etc/pam.d/vmtoolsd
/etc/init.d
/etc/init.d/open-vm-tools
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/open-vm-tools.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/60-open-vm-tools.rules
/usr/share/man/man8/mount.vmhgfs.8.gz
/usr/lib/libguestlib.so.0
/usr/lib/libvmtools.so.0
/usr/lib/libvmtools.so
/usr/lib/libguestlib.so
/usr/lib/libhgfs.so.0
/usr/lib/libhgfs.so
/sbin/vmware-hgfsmounter
