# var/lib/dpkg/info/pmacct.conffiles
/etc/pmacct/pmacctd.conf
/etc/init.d/pmacct
/etc/default/pmacct
