# var/lib/dpkg/info/pmacct.list
/.
/usr
/usr/sbin
/usr/sbin/uacctd
/usr/sbin/nfacctd
/usr/sbin/pmacctd
/usr/sbin/sfacctd
/usr/share
/usr/share/doc
/usr/share/doc/pmacct
/usr/share/doc/pmacct/copyright
/usr/share/doc/pmacct/SIGNALS
/usr/share/doc/pmacct/KNOWN-BUGS
/usr/share/doc/pmacct/PLUGINS
/usr/share/doc/pmacct/TOOLS
/usr/share/doc/pmacct/FAQS.gz
/usr/share/doc/pmacct/changelog.Debian.gz
/usr/share/doc/pmacct/CONFIG-KEYS.gz
/usr/share/doc/pmacct/README
/usr/share/doc/pmacct/UPGRADE.gz
/usr/share/doc/pmacct/TODO
/usr/share/doc/pmacct/sql
/usr/share/doc/pmacct/sql/pmacct-create-table_v7.pgsql
/usr/share/doc/pmacct/sql/pmacct-create-db_v6.mysql
/usr/share/doc/pmacct/sql/pmacct-create-table_v5.pgsql
/usr/share/doc/pmacct/sql/pmacct-create-db.pgsql
/usr/share/doc/pmacct/sql/pmacct-grant-db.mysql
/usr/share/doc/pmacct/sql/pmacct-create-table_v1.pgsql
/usr/share/doc/pmacct/sql/pmacct-create-table_v4.sqlite3
/usr/share/doc/pmacct/sql/pmacct-create-table_v6.sqlite3
/usr/share/doc/pmacct/sql/pmacct-create-db_v5.mysql
/usr/share/doc/pmacct/sql/README.64bit
/usr/share/doc/pmacct/sql/README.sqlite3.gz
/usr/share/doc/pmacct/sql/README.etype
/usr/share/doc/pmacct/sql/README.agent_id2
/usr/share/doc/pmacct/sql/pmacct-create-db_bgp_v1.mysql
/usr/share/doc/pmacct/sql/pmacct-create-db_v7.mysql
/usr/share/doc/pmacct/sql/pmacct-create-db_v3.mysql
/usr/share/doc/pmacct/sql/pmacct-create-table_bgp_v1.sqlite3
/usr/share/doc/pmacct/sql/pmacct-create-table_v7.sqlite3
/usr/share/doc/pmacct/sql/pmacct-create-table_v2.pgsql
/usr/share/doc/pmacct/sql/README.cos
/usr/share/doc/pmacct/sql/README.mask
/usr/share/doc/pmacct/sql/pmacct-create-table_v4.pgsql
/usr/share/doc/pmacct/sql/README.mpls_vpn_rd
/usr/share/doc/pmacct/sql/pmacct-create-table_v8.sqlite3
/usr/share/doc/pmacct/sql/pmacct-create-db_v1.mysql
/usr/share/doc/pmacct/sql/pmacct-create-table_v2.sqlite3
/usr/share/doc/pmacct/sql/pmacct-create-table_v3.pgsql
/usr/share/doc/pmacct/sql/README.iface
/usr/share/doc/pmacct/sql/pmacct-create-db_v8.mysql
/usr/share/doc/pmacct/sql/pmacct-create-table_v6.pgsql
/usr/share/doc/pmacct/sql/pmacct-create-table_v1.sqlite3
/usr/share/doc/pmacct/sql/pmacct-create-db_v4.mysql
/usr/share/doc/pmacct/sql/pmacct-create-table_bgp_v1.pgsql
/usr/share/doc/pmacct/sql/README.mysql.gz
/usr/share/doc/pmacct/sql/README.pgsql.gz
/usr/share/doc/pmacct/sql/README.IPv6
/usr/share/doc/pmacct/sql/pmacct-create-table_v3.sqlite3
/usr/share/doc/pmacct/sql/pmacct-create-table_v5.sqlite3
/usr/share/doc/pmacct/sql/pmacct-create-db_v2.mysql
/usr/share/doc/pmacct/QUICKSTART.gz
/usr/share/doc/pmacct/changelog.gz
/usr/share/doc/pmacct/INTERNALS.gz
/usr/share/doc/pmacct/README.Debian
/usr/bin
/usr/bin/pmacct
/var
/var/lib
/var/lib/pmacct
/etc
/etc/pmacct
/etc/pmacct/pmacctd.conf
/etc/init.d
/etc/init.d/pmacct
/etc/default
/etc/default/pmacct
