# var/lib/dpkg/info/pptpd.list
/.
/etc
/etc/ppp
/etc/ppp/pptpd-options
/etc/pptpd.conf
/etc/init.d
/etc/init.d/pptpd
/usr
/usr/lib
/usr/lib/pptpd
/usr/lib/pptpd/pptpd-logwtmp.so
/usr/sbin
/usr/sbin/pptpd
/usr/sbin/pptpctrl
/usr/share
/usr/share/doc
/usr/share/doc/pptpd
/usr/share/doc/pptpd/README.portslave
/usr/share/doc/pptpd/AUTHORS
/usr/share/doc/pptpd/changelog.gz
/usr/share/doc/pptpd/README.logwtmp
/usr/share/doc/pptpd/examples
/usr/share/doc/pptpd/examples/options.pptpd.gz
/usr/share/doc/pptpd/examples/chap-secrets
/usr/share/doc/pptpd/examples/pptpd.conf
/usr/share/doc/pptpd/README.slirp
/usr/share/doc/pptpd/poptop_ads_howto
/usr/share/doc/pptpd/poptop_ads_howto/poptop_ads_howto_9.htm
/usr/share/doc/pptpd/poptop_ads_howto/poptop_ads_howto_a2.htm
/usr/share/doc/pptpd/poptop_ads_howto/poptop_ads_howto_11.htm
/usr/share/doc/pptpd/poptop_ads_howto/poptop_ads_howto_10.htm
/usr/share/doc/pptpd/poptop_ads_howto/poptop_ads_howto_1.htm
/usr/share/doc/pptpd/poptop_ads_howto/poptop_ads_howto_2.htm
/usr/share/doc/pptpd/poptop_ads_howto/poptop_ads_howto_5.htm
/usr/share/doc/pptpd/poptop_ads_howto/poptop_ads_howto_6.htm
/usr/share/doc/pptpd/poptop_ads_howto/test.txt
/usr/share/doc/pptpd/poptop_ads_howto/poptop_ads_howto_8.htm
/usr/share/doc/pptpd/poptop_ads_howto/diagram1.jpg
/usr/share/doc/pptpd/poptop_ads_howto/poptop_ads_howto_4.htm
/usr/share/doc/pptpd/poptop_ads_howto/poptop_ads_howto_a1.htm
/usr/share/doc/pptpd/poptop_ads_howto/CVS
/usr/share/doc/pptpd/poptop_ads_howto/CVS/Repository
/usr/share/doc/pptpd/poptop_ads_howto/CVS/Entries
/usr/share/doc/pptpd/poptop_ads_howto/CVS/Root
/usr/share/doc/pptpd/poptop_ads_howto/poptop_ads_howto_3.htm
/usr/share/doc/pptpd/poptop_ads_howto/poptop_ads_howto_12.htm
/usr/share/doc/pptpd/poptop_ads_howto/poptop_ads_howto_7.htm
/usr/share/doc/pptpd/poptop_ads_howto/poptop_ads_howto_6a.htm
/usr/share/doc/pptpd/changelog.Debian.gz
/usr/share/doc/pptpd/README.inetd
/usr/share/doc/pptpd/TODO.gz
/usr/share/doc/pptpd/README
/usr/share/doc/pptpd/copyright
/usr/share/doc/pptpd/README.Debian
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/pptpd.8.gz
/usr/share/man/man8/pptpctrl.8.gz
/usr/share/man/man5
/usr/share/man/man5/pptpd.conf.5.gz
/usr/share/doc/pptpd/ChangeLog.gz
