# var/lib/dpkg/info/pptpd.conffiles
/etc/pptpd.conf
/etc/ppp/pptpd-options
/etc/init.d/pptpd
/etc/ppp/pptpd-options
/etc/pptpd.conf
/etc/init.d/pptpd
