# var/lib/dpkg/info/radvd.list
/.
/usr
/usr/sbin
/usr/sbin/radvd
/usr/sbin/radvdump
/usr/share
/usr/share/doc
/usr/share/doc/radvd
/usr/share/doc/radvd/copyright
/usr/share/doc/radvd/README
/usr/share/doc/radvd/TODO
/usr/share/doc/radvd/changelog.gz
/usr/share/doc/radvd/README.Debian
/usr/share/doc/radvd/examples
/usr/share/doc/radvd/examples/radvd.conf.example
/usr/share/doc/radvd/examples/simple-radvd.conf
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/radvdump.8.gz
/usr/share/man/man8/radvd.8.gz
/usr/share/man/man5
/usr/share/man/man5/radvd.conf.5.gz
/etc
/etc/init.d
/etc/init.d/radvd
