# var/lib/dpkg/info/snmpd.conffiles
/etc/init.d/snmpd
/etc/default/snmpd
/etc/snmp/snmptrapd.conf
/etc/snmp/snmpd.conf
