# var/lib/dpkg/info/snmpd.list
/.
/usr
/usr/sbin
/usr/sbin/snmptrapd
/usr/sbin/snmpd
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/snmpd
/usr/share/doc
/usr/share/doc/snmpd
/usr/share/doc/snmpd/copyright
/usr/share/doc/snmpd/changelog.gz
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/snmpd.8.gz
/usr/share/man/man8/snmptrapd.8.gz
/usr/share/man/man5
/usr/share/man/man5/snmptrapd.conf.5.gz
/usr/share/man/man5/snmpd.conf.5.gz
/usr/share/snmp
/usr/share/snmp/snmp_perl.pl
/usr/share/snmp/snmp_perl_trapd.pl
/usr/share/snmp/snmpconf-data
/usr/share/snmp/snmpconf-data/snmptrapd-data
/usr/share/snmp/snmpconf-data/snmptrapd-data/runtime
/usr/share/snmp/snmpconf-data/snmptrapd-data/traphandle
/usr/share/snmp/snmpconf-data/snmptrapd-data/authentication
/usr/share/snmp/snmpconf-data/snmptrapd-data/snmpconf-config
/usr/share/snmp/snmpconf-data/snmptrapd-data/logging
/usr/share/snmp/snmpconf-data/snmptrapd-data/formatting
/usr/share/snmp/snmpconf-data/snmpd-data
/usr/share/snmp/snmpconf-data/snmpd-data/operation
/usr/share/snmp/snmpconf-data/snmpd-data/system
/usr/share/snmp/snmpconf-data/snmpd-data/trapsinks
/usr/share/snmp/snmpconf-data/snmpd-data/monitor
/usr/share/snmp/snmpconf-data/snmpd-data/basic_setup
/usr/share/snmp/snmpconf-data/snmpd-data/extending
/usr/share/snmp/snmpconf-data/snmpd-data/snmpconf-config
/usr/share/snmp/snmpconf-data/snmpd-data/acl
/etc
/etc/init.d
/etc/init.d/snmpd
/etc/default
/etc/default/snmpd
/etc/snmp
/etc/snmp/snmptrapd.conf
/etc/snmp/snmpd.conf
