# var/lib/dpkg/info/squid3.list
/.
/var
/var/run
/var/spool
/var/spool/squid3
/var/log
/var/log/squid3
/etc
/etc/logrotate.d
/etc/logrotate.d/squid3
/etc/init.d
/etc/init.d/squid3
/etc/resolvconf
/etc/resolvconf/update-libc.d
/etc/resolvconf/update-libc.d/squid3
/etc/squid3
/etc/squid3/msntauth.conf
/etc/squid3/squid.conf
/etc/squid3/errorpage.css
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/squid3_pam_auth.8.gz
/usr/share/man/man8/squid3.8.gz
/usr/share/man/man8/squid3_unix_group.8.gz
/usr/share/man/man8/squid3_ldap_group.8.gz
/usr/share/man/man8/squid3_radius_auth.8.gz
/usr/share/man/man8/squid3_session.8.gz
/usr/share/man/man8/squid3_ncsa_auth.8.gz
/usr/share/man/man8/squid3_ldap_auth.8.gz
/usr/share/man/man8/squid3_db_auth.8.gz
/usr/share/doc
/usr/share/doc/squid3
/usr/share/doc/squid3/changelog.Debian.gz
/usr/share/doc/squid3/README
/usr/share/doc/squid3/changelog.gz
/usr/share/doc/squid3/copyright
/usr/share/doc/squid3/README.Debian
/usr/lib
/usr/lib/squid3
/usr/lib/squid3/yp_auth
/usr/lib/squid3/squid_ldap_auth
/usr/lib/squid3/squid_ldap_group
/usr/lib/squid3/squid_kerb_auth_test
/usr/lib/squid3/sasl_auth
/usr/lib/squid3/diskd
/usr/lib/squid3/smb_auth
/usr/lib/squid3/unlinkd
/usr/lib/squid3/ntlm_smb_lm_auth
/usr/lib/squid3/digest_pw_auth
/usr/lib/squid3/wbinfo_group.pl
/usr/lib/squid3/digest_ldap_auth
/usr/lib/squid3/pop3.pl
/usr/lib/squid3/ncsa_auth
/usr/lib/squid3/getpwname_auth
/usr/lib/squid3/squid_kerb_auth
/usr/lib/squid3/smb_auth.sh
/usr/lib/squid3/squid_db_auth
/usr/lib/squid3/negotiate_kerb_auth_test
/usr/lib/squid3/squid_radius_auth
/usr/lib/squid3/ip_user_check
/usr/lib/squid3/msnt_auth
/usr/lib/squid3/squid_session
/usr/lib/squid3/negotiate_kerb_auth
/usr/lib/squid3/squid_unix_group
/usr/lib/squid3/smb_auth.pl
/usr/lib/squid3/pam_auth
/usr/sbin
/usr/sbin/squid3
