# var/lib/dpkg/info/strongswan-starter.conffiles
/etc/ipsec.conf
/etc/ipsec.secrets
/etc/init.d/ipsec
