# var/lib/dpkg/info/vyatta-quagga.conffiles
/etc/pam.d/quagga
/etc/init.d/vyatta-quagga
