# var/lib/dpkg/info/vyatta-dhcp3-relay.list
/.
/usr
/usr/sbin
/usr/sbin/dhcrelay3
/usr/share
/usr/share/doc
/usr/share/doc/vyatta-dhcp3-relay
/usr/share/doc/vyatta-dhcp3-relay/copyright
/usr/share/doc/vyatta-dhcp3-relay/dhcp-on-linux.txt
/usr/share/doc/vyatta-dhcp3-relay/changelog.gz
/usr/share/doc/vyatta-dhcp3-relay/README.Debian
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/dhcrelay3.8.gz
/usr/share/man/man5
/var
/var/run
/etc
/etc/init.d
/etc/init.d/vyatta-dhcp3-relay
