# var/lib/dpkg/info/vyatta-keepalived.list
/.
/usr
/usr/sbin
/usr/sbin/keepalived
/usr/share
/usr/share/doc
/usr/share/doc/vyatta-keepalived
/usr/share/doc/vyatta-keepalived/copyright
/usr/share/doc/vyatta-keepalived/keepalived.conf.SYNOPSIS.gz
/usr/share/doc/vyatta-keepalived/NEWS.gz
/usr/share/doc/vyatta-keepalived/changelog.Debian.gz
/usr/share/doc/vyatta-keepalived/README
/usr/share/doc/vyatta-keepalived/TODO
/usr/share/doc/vyatta-keepalived/AUTHOR
/usr/share/doc/vyatta-keepalived/changelog.gz
/usr/share/doc/vyatta-keepalived/CONTRIBUTORS
/usr/share/doc/vyatta-keepalived/samples
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.SSL_GET
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.SMTP_CHECK
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.inhibit
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.vrrp.static_ipaddress
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.virtualhost
/usr/share/doc/vyatta-keepalived/samples/dh1024.pem
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.fwmark
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.virtual_server_group
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.vrrp.sync
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.misc_check
/usr/share/doc/vyatta-keepalived/samples/client.pem
/usr/share/doc/vyatta-keepalived/samples/root.pem
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.vrrp.localcheck
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.status_code
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.misc_check_arg
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.vrrp.routes
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.sample
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.track_interface
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.IPv6
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.vrrp
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.quorum
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.vrrp.scripts
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.vrrp.lvs_syncd
/usr/share/doc/vyatta-keepalived/samples/sample.misccheck.smbcheck.sh
/usr/share/doc/vyatta-keepalived/samples/keepalived.conf.HTTP_GET.port
/usr/share/doc/vyatta-keepalived/README.Debian
/usr/share/doc/keepalived
/usr/share/doc/keepalived/samples
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/genhash.1.gz
/usr/share/man/man8
/usr/share/man/man8/keepalived.8.gz
/usr/share/man/man5
/usr/share/man/man5/keepalived.conf.5.gz
/usr/share/snmp
/usr/share/snmp/mibs
/usr/share/snmp/mibs/KEEPALIVED-MIB
/usr/share/snmp/mibs/VRRP-MIB
/usr/bin
/usr/bin/genhash
/etc
/etc/keepalived
/etc/keepalived/keepalived.conf
/etc/keepalived/samples
/etc/keepalived/samples/keepalived.conf.SSL_GET
/etc/keepalived/samples/keepalived.conf.SMTP_CHECK
/etc/keepalived/samples/keepalived.conf.inhibit
/etc/keepalived/samples/keepalived.conf.vrrp.static_ipaddress
/etc/keepalived/samples/keepalived.conf.virtualhost
/etc/keepalived/samples/dh1024.pem
/etc/keepalived/samples/keepalived.conf.fwmark
/etc/keepalived/samples/keepalived.conf.virtual_server_group
/etc/keepalived/samples/keepalived.conf.vrrp.sync
/etc/keepalived/samples/keepalived.conf.misc_check
/etc/keepalived/samples/client.pem
/etc/keepalived/samples/root.pem
/etc/keepalived/samples/keepalived.conf.vrrp.localcheck
/etc/keepalived/samples/keepalived.conf.status_code
/etc/keepalived/samples/keepalived.conf.misc_check_arg
/etc/keepalived/samples/keepalived.conf.vrrp.routes
/etc/keepalived/samples/keepalived.conf.sample
/etc/keepalived/samples/keepalived.conf.track_interface
/etc/keepalived/samples/keepalived.conf.IPv6
/etc/keepalived/samples/keepalived.conf.vrrp
/etc/keepalived/samples/keepalived.conf.quorum
/etc/keepalived/samples/keepalived.conf.vrrp.scripts
/etc/keepalived/samples/keepalived.conf.vrrp.lvs_syncd
/etc/keepalived/samples/sample.misccheck.smbcheck.sh
/etc/keepalived/samples/keepalived.conf.HTTP_GET.port
/etc/sysconfig
/etc/sysconfig/keepalived
/etc/rc.d
/etc/rc.d/init.d
/etc/rc.d/init.d/keepalived
/etc/init.d
/etc/init.d/vyatta-keepalived
