# var/lib/dpkg/info/vyatta-wanloadbalance.conffiles
/etc/ppp/ip-up.d/ppp-wanlb
/etc/dhcp3/dhclient-exit-hooks.d/dhcp-wanlb
/etc/loadbalance.conf
/etc/init.d/vyatta-wanloadbalance
