# var/lib/dpkg/info/vyatta-openvpn.postinst
#!/bin/bash

prefix=/opt/vyatta
exec_prefix=${prefix}
sysconfdir=${prefix}/etc
bindir=${exec_prefix}/bin
sbindir=${exec_prefix}/sbin

# remove init of daemons that we start/stop
for init in openvpn; do
  update-rc.d -f ${init} remove >/dev/null
done

mkdir -p /opt/vyatta/etc/openvpn/ccd
mkdir -p /opt/vyatta/etc/openvpn/status

