# var/lib/dpkg/info/vyatta-conntrack-sync.list
/.
/usr
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/vyatta-conntrack-sync
/usr/share/doc
/usr/share/doc/vyatta-conntrack-sync
/usr/share/doc/vyatta-conntrack-sync/copyright
/usr/share/doc/vyatta-conntrack-sync/NEWS.gz
/usr/share/doc/vyatta-conntrack-sync/changelog.Debian.gz
/usr/share/doc/vyatta-conntrack-sync/README
/usr/share/doc/vyatta-conntrack-sync/changelog.gz
/opt
/opt/vyatta
/opt/vyatta/sbin
/opt/vyatta/sbin/vyatta-conntrack-sync.pl
/opt/vyatta/sbin/vyatta-vrrp-conntracksync.sh
/opt/vyatta/share
/opt/vyatta/share/vyatta-cfg
/opt/vyatta/share/vyatta-cfg/templates
/opt/vyatta/share/vyatta-cfg/templates/service
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/expect-sync
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/expect-sync/node.def
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/accept-protocol
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/accept-protocol/node.def
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/node.def
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/ignore-address
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/ignore-address/node.def
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/ignore-address/ipv4
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/ignore-address/ipv4/node.def
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/sync-queue-size
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/sync-queue-size/node.def
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/failover-mechanism
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/failover-mechanism/node.def
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/failover-mechanism/vrrp
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/failover-mechanism/vrrp/node.def
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/failover-mechanism/vrrp/sync-group
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/failover-mechanism/vrrp/sync-group/node.def
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/failover-mechanism/cluster
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/failover-mechanism/cluster/node.def
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/failover-mechanism/cluster/group
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/failover-mechanism/cluster/group/node.def
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/interface
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/interface/node.def
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/mcast-group
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/mcast-group/node.def
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/event-listen-queue-size
/opt/vyatta/share/vyatta-cfg/templates/service/conntrack-sync/event-listen-queue-size/node.def
/opt/vyatta/share/perl5
/opt/vyatta/share/perl5/Vyatta
/opt/vyatta/share/perl5/Vyatta/ConntrackSync.pm
/opt/vyatta/share/vyatta-op
/opt/vyatta/share/vyatta-op/templates
/opt/vyatta/share/vyatta-op/templates/reset
/opt/vyatta/share/vyatta-op/templates/reset/conntrack-sync
/opt/vyatta/share/vyatta-op/templates/reset/conntrack-sync/node.def
/opt/vyatta/share/vyatta-op/templates/reset/conntrack-sync/external-cache
/opt/vyatta/share/vyatta-op/templates/reset/conntrack-sync/external-cache/node.def
/opt/vyatta/share/vyatta-op/templates/reset/conntrack-sync/internal-cache
/opt/vyatta/share/vyatta-op/templates/reset/conntrack-sync/internal-cache/node.def
/opt/vyatta/share/vyatta-op/templates/restart
/opt/vyatta/share/vyatta-op/templates/restart/conntrack-sync
/opt/vyatta/share/vyatta-op/templates/restart/conntrack-sync/node.def
/opt/vyatta/share/vyatta-op/templates/show
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/statistics
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/statistics/node.def
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/statistics/main
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/statistics/main/node.def
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/statistics/expect
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/statistics/expect/node.def
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/node.def
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/external-cache
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/external-cache/node.def
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/external-cache/main
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/external-cache/main/node.def
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/external-cache/expect
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/external-cache/expect/node.def
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/internal-cache
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/internal-cache/node.def
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/internal-cache/main
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/internal-cache/main/node.def
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/internal-cache/expect
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/internal-cache/expect/node.def
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/status
/opt/vyatta/share/vyatta-op/templates/show/conntrack-sync/status/node.def
/opt/vyatta/etc
/opt/vyatta/etc/config-migrate
/opt/vyatta/etc/config-migrate/current
/opt/vyatta/etc/config-migrate/current/conntrack-sync@1
/opt/vyatta/bin
/opt/vyatta/bin/sudo-users
/opt/vyatta/bin/sudo-users/vyatta-op-conntrack-sync.pl
/opt/vyatta/bin/sudo-users/vyatta-show-ctsync-caches.pl
/etc
/etc/init.d
/etc/init.d/vyatta-cluster-conntracksync
