# var/lib/dpkg/info/vyatta-cfg-op-pppoe.list
/.
/opt
/opt/vyatta
/opt/vyatta/share
/opt/vyatta/share/vyatta-cfg
/opt/vyatta/share/vyatta-cfg/templates
/opt/vyatta/share/vyatta-cfg/templates/interfaces
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/connect-on-demand
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/connect-on-demand/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/idle-timeout
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/idle-timeout/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/local-address
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/local-address/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/disable
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/disable/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/access-concentrator
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/access-concentrator/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/password
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/password/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/user-id
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/user-id/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/service-name
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/service-name/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/mtu
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/mtu/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/default-route
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/default-route/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/remote-address
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/remote-address/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/name-server
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.tag/name-server/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/pppoe/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/connect-on-demand
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/connect-on-demand/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/idle-timeout
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/idle-timeout/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/local-address
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/local-address/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/access-concentrator
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/access-concentrator/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/password
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/password/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/user-id
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/user-id/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/service-name
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/service-name/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/mtu
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/mtu/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/default-route
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/default-route/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/remote-address
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/remote-address/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/name-server
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.tag/name-server/node.def
/opt/vyatta/share/vyatta-cfg/templates/interfaces/ethernet/node.tag/vif/node.tag/pppoe/node.def
/opt/vyatta/share/vyatta-op
/opt/vyatta/share/vyatta-op/templates
/opt/vyatta/share/vyatta-op/templates/show
/opt/vyatta/share/vyatta-op/templates/show/interfaces
/opt/vyatta/share/vyatta-op/templates/show/interfaces/pppoe
/opt/vyatta/share/vyatta-op/templates/show/interfaces/pppoe/node.tag
/opt/vyatta/share/vyatta-op/templates/show/interfaces/pppoe/node.tag/node.def
/opt/vyatta/share/vyatta-op/templates/show/interfaces/pppoe/node.tag/log
/opt/vyatta/share/vyatta-op/templates/show/interfaces/pppoe/node.tag/log/tail
/opt/vyatta/share/vyatta-op/templates/show/interfaces/pppoe/node.tag/log/tail/node.def
/opt/vyatta/share/vyatta-op/templates/show/interfaces/pppoe/node.tag/log/node.def
/opt/vyatta/share/vyatta-op/templates/show/interfaces/pppoe/node.def
/opt/vyatta/share/vyatta-op/templates/monitor
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/unlimited
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/unlimited/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/unlimited/filter
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/unlimited/filter/node.tag
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/unlimited/filter/node.tag/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/unlimited/filter/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/filter
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/filter/node.tag
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/filter/node.tag/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/filter/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/detail
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/detail/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/detail/unlimited
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/detail/unlimited/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/detail/unlimited/filter
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/detail/unlimited/filter/node.tag
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/detail/unlimited/filter/node.tag/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/detail/unlimited/filter/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/detail/filter
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/detail/filter/node.tag
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/detail/filter/node.tag/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/detail/filter/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/save
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/save/node.tag
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/save/node.tag/size
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/save/node.tag/size/node.tag
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/save/node.tag/size/node.tag/files
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/save/node.tag/size/node.tag/files/node.tag
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/save/node.tag/size/node.tag/files/node.tag/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/save/node.tag/size/node.tag/files/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/save/node.tag/size/node.tag/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/save/node.tag/size/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/save/node.tag/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.tag/traffic/save/node.def
/opt/vyatta/share/vyatta-op/templates/monitor/interfaces/pppoe/node.def
/opt/vyatta/etc
/opt/vyatta/etc/pppoe-provider-template
/usr
/usr/share
/usr/share/doc
/usr/share/doc/vyatta-cfg-op-pppoe
/usr/share/doc/vyatta-cfg-op-pppoe/changelog.gz
/usr/share/doc/vyatta-cfg-op-pppoe/NEWS.gz
/usr/share/doc/vyatta-cfg-op-pppoe/README
/usr/share/doc/vyatta-cfg-op-pppoe/copyright
/usr/share/doc/vyatta-cfg-op-pppoe/changelog.Debian.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/vyatta-cfg-op-pppoe
/usr/share/linda
/usr/share/linda/overrides
/usr/share/linda/overrides/vyatta-cfg-op-pppoe
/etc
/etc/init.d
/etc/init.d/vyatta-pppoe
/etc/ppp
/etc/ppp/ip-pre-up.d
/etc/ppp/ip-pre-up.d/0002rename-pppoe
