# var/lib/dpkg/info/vyatta-webgui.list
/.
/usr
/usr/sbin
/usr/sbin/chunker_cmd
/usr/sbin/chunker
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/vyatta-webgui
/usr/share/doc
/usr/share/doc/vyatta-webgui
/usr/share/doc/vyatta-webgui/copyright
/usr/share/doc/vyatta-webgui/NEWS.gz
/usr/share/doc/vyatta-webgui/changelog.Debian.gz
/usr/share/doc/vyatta-webgui/README
/usr/share/doc/vyatta-webgui/changelog.gz
/usr/lib
/usr/lib/cgi-bin
/usr/lib/cgi-bin/webgui-wrap
/usr/lib/cgi-bin/webgui
/opt
/opt/vyatta
/opt/vyatta/tmp
/opt/vyatta/tmp/webgui
/opt/vyatta/tmp/webgui/dummy
/opt/vyatta/var
/opt/vyatta/var/run
/opt/vyatta/sbin
/opt/vyatta/sbin/vyatta-update-webgui-listen-addr.pl
/opt/vyatta/share
/opt/vyatta/share/vyatta-cfg
/opt/vyatta/share/vyatta-cfg/templates
/opt/vyatta/share/vyatta-cfg/templates/service
/opt/vyatta/share/vyatta-cfg/templates/service/https
/opt/vyatta/share/vyatta-cfg/templates/service/https/http-redirect
/opt/vyatta/share/vyatta-cfg/templates/service/https/http-redirect/node.def
/opt/vyatta/share/vyatta-cfg/templates/service/https/node.def
/opt/vyatta/share/vyatta-cfg/templates/service/https/listen-address
/opt/vyatta/share/vyatta-cfg/templates/service/https/listen-address/node.def
/opt/vyatta/share/vyatta-op
/opt/vyatta/share/vyatta-op/templates
/opt/vyatta/share/vyatta-op/templates/restart
/opt/vyatta/share/vyatta-op/templates/restart/https
/opt/vyatta/share/vyatta-op/templates/restart/https/node.def
/var
/var/www
/var/www/index.html
/var/www/vyos
/var/www/vyos/images
/var/www/vyos/images/vy_icon.png
/var/www/vyos/images/vyos_logo_full.png
/var/www/vyos/main.html
/etc
/etc/init.d
/etc/init.d/vyatta-webgui-chunker
