# var/lib/dpkg/info/vyatta-cfg-op-pppoe.conffiles
/etc/init.d/vyatta-pppoe
/etc/ppp/ip-pre-up.d/0002rename-pppoe
