# var/lib/dpkg/info/vyatta-dhcp3-server.conffiles
/etc/dhcp3/dhcpd.conf
/etc/init.d/vyatta-dhcp3-server
