# var/lib/dpkg/info/vyatta-conntrack-sync.conffiles
/etc/init.d/vyatta-cluster-conntracksync
