# var/lib/dpkg/info/vyatta-cfg.list
/.
/usr
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/vyatta-cfg
/usr/share/doc
/usr/share/doc/vyatta-cfg
/usr/share/doc/vyatta-cfg/copyright
/usr/share/doc/vyatta-cfg/NEWS.gz
/usr/share/doc/vyatta-cfg/README
/usr/share/doc/vyatta-cfg/changelog.gz
/opt
/opt/vyatta
/opt/vyatta/sbin
/opt/vyatta/sbin/vyatta-validate-type.pl
/opt/vyatta/sbin/vyatta-boot-config-loader
/opt/vyatta/sbin/priority.pl
/opt/vyatta/sbin/vyatta-load-config.pl
/opt/vyatta/sbin/vyatta-strip-migration-comments.pl
/opt/vyatta/sbin/priority
/opt/vyatta/sbin/check_tmpl
/opt/vyatta/sbin/vyatta-cli-expand-var.pl
/opt/vyatta/sbin/exe_action
/opt/vyatta/sbin/my_cli_bin
/opt/vyatta/sbin/my_cli_shell_api
/opt/vyatta/sbin/vyatta-config-gen-sets.pl
/opt/vyatta/sbin/vyatta-find-type.pl
/opt/vyatta/sbin/vyatta-save-config.pl
/opt/vyatta/sbin/dump
/opt/vyatta/sbin/vyatta-log-commit.pl
/opt/vyatta/sbin/vyatta-cfg-cmd-wrapper
/opt/vyatta/sbin/vyatta-cfg-notify
/opt/vyatta/sbin/vyatta-output-config.pl
/opt/vyatta/sbin/vyatta-check-typeless-node.pl
/opt/vyatta/share
/opt/vyatta/share/enumeration
/opt/vyatta/share/enumeration/existing-interfaces
/opt/vyatta/share/vyatta-cfg
/opt/vyatta/share/vyatta-cfg/functions
/opt/vyatta/share/vyatta-cfg/functions/interpreter
/opt/vyatta/share/vyatta-cfg/functions/interpreter/vyatta-cfg-run
/opt/vyatta/etc
/opt/vyatta/etc/functions
/opt/vyatta/etc/functions/script-template
/opt/vyatta/etc/shell
/opt/vyatta/etc/shell/level
/opt/vyatta/etc/shell/level/users
/opt/vyatta/etc/shell/level/users/allowed-cfg
/opt/vyatta/etc/shell/level/users/restricted-mode
/opt/vyatta/etc/shell/level/users/allowed-pipe
/opt/vyatta/etc/shell/level/admin
/opt/vyatta/etc/shell/level/admin/restricted-mode
/etc
/etc/bash_completion.d
/etc/bash_completion.d/vyatta-cfg
/etc/logrotate.d
/etc/logrotate.d/vyatta-config-logs
/etc/dhcp3
/etc/dhcp3/dhclient-exit-hooks.d
/etc/dhcp3/dhclient-exit-hooks.d/vyatta-dhclient-hook
/etc/init.d
/etc/init.d/vyatta-router
/etc/default
/etc/default/vyatta-load-boot
/etc/default/vyatta-cfg
/bin
/opt/vyatta/sbin/my_set
/opt/vyatta/sbin/my_rename
/opt/vyatta/sbin/my_commit
/opt/vyatta/sbin/my_move
/opt/vyatta/sbin/my_comment
/opt/vyatta/sbin/my_delete
/opt/vyatta/sbin/my_copy
/opt/vyatta/sbin/my_discard
/bin/cli-shell-api
