# var/lib/dpkg/info/vyatta-quagga.list
/.
/usr
/usr/include
/usr/include/quagga
/usr/include/quagga/ospfd
/usr/include/quagga/ospfd/ospf_asbr.h
/usr/include/quagga/ospfd/ospf_opaque.h
/usr/include/quagga/ospfd/ospf_nsm.h
/usr/include/quagga/ospfd/ospf_lsa.h
/usr/include/quagga/ospfd/ospf_lsdb.h
/usr/include/quagga/ospfd/ospf_ism.h
/usr/include/quagga/ospfd/ospf_api.h
/usr/include/quagga/ospfd/ospfd.h
/usr/include/quagga/ospfd/ospf_dump.h
/usr/include/quagga/network.h
/usr/include/quagga/hash.h
/usr/include/quagga/pqueue.h
/usr/include/quagga/zebra.h
/usr/include/quagga/linklist.h
/usr/include/quagga/jhash.h
/usr/include/quagga/workqueue.h
/usr/include/quagga/version.h
/usr/include/quagga/ospfapi
/usr/include/quagga/ospfapi/ospf_apiclient.h
/usr/include/quagga/plist.h
/usr/include/quagga/memory.h
/usr/include/quagga/log.h
/usr/include/quagga/sockopt.h
/usr/include/quagga/prefix.h
/usr/include/quagga/keychain.h
/usr/include/quagga/routemap.h
/usr/include/quagga/getopt.h
/usr/include/quagga/distribute.h
/usr/include/quagga/if_rmap.h
/usr/include/quagga/sigevent.h
/usr/include/quagga/str.h
/usr/include/quagga/if.h
/usr/include/quagga/vector.h
/usr/include/quagga/privs.h
/usr/include/quagga/table.h
/usr/include/quagga/vty.h
/usr/include/quagga/thread.h
/usr/include/quagga/stream.h
/usr/include/quagga/route_types.h
/usr/include/quagga/zclient.h
/usr/include/quagga/memtypes.h
/usr/include/quagga/checksum.h
/usr/include/quagga/sockunion.h
/usr/include/quagga/buffer.h
/usr/include/quagga/filter.h
/usr/include/quagga/zassert.h
/usr/include/quagga/smux.h
/usr/include/quagga/command.h
/usr/include/quagga/md5.h
/usr/sbin
/usr/sbin/ospfd
/usr/sbin/zebra
/usr/sbin/bgpd
/usr/sbin/ripd
/usr/sbin/watchquagga
/usr/sbin/ospfclient
/usr/sbin/ripngd
/usr/sbin/ospf6d
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/vyatta-quagga
/usr/share/doc
/usr/share/doc/quagga
/usr/share/doc/quagga/examples
/usr/share/doc/quagga/examples/vtysh.conf.sample
/usr/share/doc/quagga/examples/bgpd.conf.sample2
/usr/share/doc/quagga/examples/zebra.conf.sample
/usr/share/doc/quagga/examples/ospf6d.conf.sample
/usr/share/doc/quagga/examples/bgpd.conf.sample
/usr/share/doc/quagga/examples/ripd.conf.sample
/usr/share/doc/quagga/examples/ospfd.conf.sample
/usr/share/doc/quagga/examples/ripngd.conf.sample
/usr/share/doc/vyatta-quagga
/usr/share/doc/vyatta-quagga/copyright
/usr/share/doc/vyatta-quagga/NEWS.gz
/usr/share/doc/vyatta-quagga/changelog.Debian.gz
/usr/share/doc/vyatta-quagga/README
/usr/share/doc/vyatta-quagga/AUTHORS
/usr/share/doc/vyatta-quagga/changelog.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/vtysh.1.gz
/usr/share/man/man8
/usr/share/man/man8/ripngd.8.gz
/usr/share/man/man8/ospfd.8.gz
/usr/share/man/man8/bgpd.8.gz
/usr/share/man/man8/ospf6d.8.gz
/usr/share/man/man8/ripd.8.gz
/usr/share/man/man8/zebra.8.gz
/usr/share/man/man8/isisd.8.gz
/usr/share/snmp
/usr/share/snmp/mibs
/usr/share/snmp/mibs/GNOME-PRODUCT-ZEBRA-MIB.txt
/usr/lib
/usr/lib/quagga
/usr/lib/quagga/libospfapiclient.a
/usr/lib/quagga/libospf.a
/usr/lib/quagga/libospfapiclient.so.0.0.0
/usr/lib/quagga/libospf.so.0.0.0
/usr/lib/quagga/libzebra.la
/usr/lib/quagga/libzebra.a
/usr/lib/quagga/libzebra.so.0.0.0
/usr/lib/quagga/libospfapiclient.la
/usr/lib/quagga/libospf.la
/usr/bin
/usr/bin/vtysh
/var
/var/log
/var/log/quagga
/etc
/etc/quagga
/etc/pam.d
/etc/pam.d/quagga
/etc/logrotate.d
/etc/init.d
/etc/init.d/vyatta-quagga
/usr/lib/quagga/libospf.so.0
/usr/lib/quagga/libospfapiclient.so
/usr/lib/quagga/libzebra.so
/usr/lib/quagga/libospf.so
/usr/lib/quagga/libospfapiclient.so.0
/usr/lib/quagga/libzebra.so.0
