# var/lib/dpkg/info/vyatta-dhcp3-server.list
/.
/usr
/usr/sbin
/usr/sbin/dhcpd3
/usr/share
/usr/share/doc
/usr/share/doc/vyatta-dhcp3-server
/usr/share/doc/vyatta-dhcp3-server/NEWS.Debian.gz
/usr/share/doc/vyatta-dhcp3-server/copyright
/usr/share/doc/vyatta-dhcp3-server/dhcp-on-linux.txt
/usr/share/doc/vyatta-dhcp3-server/changelog.gz
/usr/share/doc/vyatta-dhcp3-server/README.Debian
/usr/share/doc/vyatta-dhcp3-server/examples
/usr/share/doc/vyatta-dhcp3-server/examples/dhcpd.conf
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/dhcpd3.8.gz
/usr/share/man/man5
/usr/share/man/man5/dhcpd.conf.5.gz
/usr/share/man/man5/dhcpd.leases.5.gz
/var
/var/run
/var/lib
/var/lib/dhcp3
/etc
/etc/dhcp3
/etc/dhcp3/dhcpd.conf
/etc/init.d
/etc/init.d/vyatta-dhcp3-server
