# var/lib/dpkg/info/vyos-nhrp.list
/.
/opt
/opt/vyatta
/opt/vyatta/sbin
/opt/vyatta/sbin/vyos-update-nhrp.pl
/opt/vyatta/share
/opt/vyatta/share/vyatta-cfg
/opt/vyatta/share/vyatta-cfg/templates
/opt/vyatta/share/vyatta-cfg/templates/protocols
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/node.def
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/redirect
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/redirect/node.def
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/map
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/map/node.tag
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/map/node.tag/cisco
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/map/node.tag/cisco/node.def
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/map/node.tag/register
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/map/node.tag/register/node.def
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/map/node.tag/nbma-address
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/map/node.tag/nbma-address/node.def
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/map/node.def
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/shortcut-target
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/shortcut-target/node.tag
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/shortcut-target/node.tag/holding-time
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/shortcut-target/node.tag/holding-time/node.def
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/shortcut-target/node.def
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/dynamic-map
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/dynamic-map/node.tag
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/dynamic-map/node.tag/nbma-domain-name
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/dynamic-map/node.tag/nbma-domain-name/node.def
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/dynamic-map/node.def
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/multicast
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/multicast/node.def
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/non-caching
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/non-caching/node.def
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/shortcut-destination
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/shortcut-destination/node.def
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/shortcut
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/shortcut/node.def
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/cisco-authentication
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/cisco-authentication/node.def
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/holding-time
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.tag/holding-time/node.def
/opt/vyatta/share/vyatta-cfg/templates/protocols/nhrp/tunnel/node.def
/opt/vyatta/share/vyatta-op
/opt/vyatta/share/vyatta-op/templates
/opt/vyatta/share/vyatta-op/templates/show
/opt/vyatta/share/vyatta-op/templates/show/nhrp
/opt/vyatta/share/vyatta-op/templates/show/nhrp/interface
/opt/vyatta/share/vyatta-op/templates/show/nhrp/interface/node.def
/opt/vyatta/share/vyatta-op/templates/show/nhrp/node.def
/opt/vyatta/share/vyatta-op/templates/show/nhrp/tunnel
/opt/vyatta/share/vyatta-op/templates/show/nhrp/tunnel/node.def
/opt/vyatta/share/vyatta-op/templates/reset
/opt/vyatta/share/vyatta-op/templates/reset/nhrp
/opt/vyatta/share/vyatta-op/templates/reset/nhrp/node.def
/opt/vyatta/share/vyatta-op/templates/reset/nhrp/flush
/opt/vyatta/share/vyatta-op/templates/reset/nhrp/flush/node.def
/opt/vyatta/share/vyatta-op/templates/reset/nhrp/flush/tunnel
/opt/vyatta/share/vyatta-op/templates/reset/nhrp/flush/tunnel/node.tag
/opt/vyatta/share/vyatta-op/templates/reset/nhrp/flush/tunnel/node.tag/node.def
/opt/vyatta/share/vyatta-op/templates/reset/nhrp/flush/tunnel/node.def
/opt/vyatta/share/vyatta-op/templates/reset/nhrp/purge
/opt/vyatta/share/vyatta-op/templates/reset/nhrp/purge/node.def
/opt/vyatta/share/vyatta-op/templates/reset/nhrp/purge/tunnel
/opt/vyatta/share/vyatta-op/templates/reset/nhrp/purge/tunnel/node.tag
/opt/vyatta/share/vyatta-op/templates/reset/nhrp/purge/tunnel/node.tag/node.def
/opt/vyatta/share/vyatta-op/templates/reset/nhrp/purge/tunnel/node.def
/opt/vyatta/etc
/opt/vyatta/etc/config-migrate
/opt/vyatta/etc/config-migrate/current
/opt/vyatta/bin
/opt/vyatta/bin/sudo-users
/usr
/usr/share
/usr/share/doc
/usr/share/doc/vyos-nhrp
/usr/share/doc/vyos-nhrp/changelog.gz
/usr/share/doc/vyos-nhrp/NEWS.gz
/usr/share/doc/vyos-nhrp/README
/usr/share/doc/vyos-nhrp/copyright
/usr/share/doc/vyos-nhrp/changelog.Debian.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/vyos-nhrp
/etc
/etc/init.d
/etc/init.d/opennhrp.init
