# var/lib/dpkg/info/xl2tpd.list
/.
/etc
/etc/init.d
/etc/init.d/xl2tpd
/etc/xl2tpd
/etc/xl2tpd/xl2tpd.conf
/etc/xl2tpd/l2tp-secrets
/etc/default
/etc/default/xl2tpd
/usr
/usr/bin
/usr/bin/pfc
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/xl2tpd.8.gz
/usr/share/man/man5
/usr/share/man/man5/l2tp-secrets.5.gz
/usr/share/man/man5/xl2tpd.conf.5.gz
/usr/share/man/man1
/usr/share/man/man1/pfc.1.gz
/usr/share/doc
/usr/share/doc/xl2tpd
/usr/share/doc/xl2tpd/CREDITS
/usr/share/doc/xl2tpd/changelog.Debian.gz
/usr/share/doc/xl2tpd/copyright
/usr/share/doc/xl2tpd/examples
/usr/share/doc/xl2tpd/examples/xl2tpd.conf
/usr/share/doc/xl2tpd/examples/xl2tpd-L2TP-CERT.conf
/usr/share/doc/xl2tpd/examples/README
/usr/share/doc/xl2tpd/examples/xl2tpd-L2TP-PSK.conf
/usr/share/doc/xl2tpd/examples/xl2tpd-L2TP-PSK-orgWIN2KXP.conf
/usr/share/doc/xl2tpd/examples/xl2tpd-L2TP-CERT-orgWIN2KXP.conf
/usr/share/doc/xl2tpd/examples/chapsecrets.sample
/usr/share/doc/xl2tpd/examples/ppp-options.xl2tpd
/usr/share/doc/xl2tpd/pfc.README
/usr/share/doc/xl2tpd/pppol2tp-2.6.23.README
/usr/share/doc/xl2tpd/changelog.gz
/usr/share/doc/xl2tpd/README.xl2tpd
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/xl2tpd
/usr/sbin
/usr/sbin/xl2tpd
